# deploy/systemd/quetta-bench-gpu-orphan-cleaner.timer
[Unit]
Description=Reclaim QuettaBench orphan GPU cleanup candidates every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=quetta-bench-gpu-orphan-cleaner.service

[Install]
WantedBy=timers.target
