# deploy/systemd/quetta-bench-orchestrator.timer
[Unit]
Description=Run QuettaBench GPU benchmark orchestration every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
RandomizedDelaySec=60s
Persistent=true
Unit=quetta-bench-orchestrator.service

[Install]
WantedBy=timers.target
